# private_dot_config/systemd/user/rclone-onedrive.service
[Unit]
Description=Mount OneDrive with rclone
Documentation=https://rclone.org/commands/rclone_mount/
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStartPre=/usr/bin/mkdir -p %h/OneDrive %h/.cache/rclone
ExecStart=/usr/bin/rclone mount OneDrive: %h/OneDrive --vfs-cache-mode full --vfs-cache-max-size 10G --vfs-cache-max-age 168h --dir-cache-time 72h --poll-interval 1m --cache-dir %h/.cache/rclone --log-level INFO --log-file %h/.cache/rclone/onedrive.log
ExecStop=/usr/bin/fusermount3 -uz %h/OneDrive
Restart=on-failure
RestartSec=10
TimeoutStopSec=20

[Install]
WantedBy=default.target
